# rootfs/etc/s6-overlay/s6-rc.d/init-ical-icloud-sync/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# run initial IP update
exec \
    s6-setuidgid abc /app/ical-icloud-sync.sh
